# openpine-gateway.service
[Unit]
Description=OpenPine Gateway — REST API + WebSocket for Pine trading system
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openpine
Group=openpine
WorkingDirectory=/opt/openpine
Environment=OPENPINE_HOME=/var/lib/openpine
Environment=OPENPINE_ROOT=/opt/openpine
Environment=OPENPINE_PYTHON=/opt/openpine/.venv/bin/python
EnvironmentFile=-/opt/openpine/.env
ExecStart=/opt/openpine/scripts/run_gateway_systemd.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openpine-gateway

[Install]
WantedBy=multi-user.target
