# openpine-gateway.service
[Unit]
Description=OpenPine Gateway — REST API + WebSocket for Pine trading system
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=moltbot1
Group=moltbot1
WorkingDirectory=[local-home]/[workspace-root]/repo-root
Environment=OPENPINE_HOME=[local-home]/[workspace-root]/repo-root
ExecStart=/usr/bin/python3 -c "from gateway.server import create_app; import uvicorn; uvicorn.run(create_app(), host='0.0.0.0', port=8080, log_level='info')"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openpine-gateway

[Install]
WantedBy=multi-user.target
